# nextion.service
[Unit]
Description = Nextion Service for "DStarRepeater" by Yosh Todo/JE3HCZ
After = network.target

[Service]
#Type = forking
Type = simple
ExecStart = /usr/local/bin/nextion

[Install]
WantedBy = multi-user.target
